# files/containers/radarr.container
[Unit]
Description=Radarr container
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=radarr
Image=lscr.io/linuxserver/radarr:latest
AutoUpdate=registry
#
PublishPort=7878:7878
#
Volume=/data/data/:/data
Volume=radarr-config:/config
#
Environment=PUID=5000
Environment=PGID=2000
Environment=UMASK=002
Environment=TZ=Europe/London

[Install]
WantedBy=multi-user.target
