# zabbix/files/zabbix-java-gateway.container
[Unit]
Description=Zabbix java-gateway container
After=zabbix-pod.service
Wants=zabbix-pod.service

[Container]
ContainerName=zabbix-java-gateway
Image=docker.io/zabbix/zabbix-java-gateway:latest
PodmanArgs=--pod=zabbix

[Install]
WantedBy=multi-user.target
